# main/.config/systemd/user/bat.service
[Unit]
Description=Bat

[Service]
ExecStart=/home/wpovell/bin/batterynotify