# main/.config/systemd/user/caffinate.service
[Unit]
Description=Caffinate

[Service]
ExecStart=/home/wpovell/bin/caffinate