# main/.config/systemd/user/caffinate.timer
[Unit]
Description=Runs caffinate every 3min

[Timer]
OnBootSec=3min
OnUnitActiveSec=3min

[Install]
WantedBy=caffinate.target
