# main/.config/systemd/user/caffinate.service
[Unit]
Description=Caffinate

[Service]
Type=Simple
ExecStart=/home/wpovell/bin/caffinate

[Install]
WantedBy=default.target